# mixcloud.timer
[Unit]
Description=Mixcloud timing

[Timer]
OnBootSec=5m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target